# deploy/linux/kanban.service
# 甲骨易经营看板 · systemd 单元（Ubuntu 26.04 · 任务书50）
# 安装见 docs/Ubuntu部署手册.md；模板内路径按实际部署目录改。
# 生产：配合 nginx-kanban.conf，本服务仅绑 127.0.0.1:8018。
#
# 语义：
#   - Restart=always + RestartSec=3：进程退出（含一键更新码 42）后自动再起
#   - ExecStart 走 start_with_rollback.sh：坏版本 (.update_rollback) 自动 git reset 一次
#   - StartLimitBurst=5 / StartLimitIntervalSec=120：短时反复崩则停下报警（对齐看门狗 5 次）
#     注意：StartLimit* 必须在 [Unit]（Ubuntu 26 / systemd 新版放 [Service] 会被忽略）
#
# 安装示例（部署目录=/opt/kanban/看板正式程序）：
#   生产机实际 User/Group=lee（数据目录属主）；模板默认也用 lee，勿盲改 kanban 除非先 chown。
#   sudo cp deploy/linux/kanban.service /etc/systemd/system/kanban.service
#   sudo systemctl daemon-reload && sudo systemctl enable --now kanban
#   journalctl -u kanban -f

[Unit]
Description=甲骨易经营看板 (kanban)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=lee
Group=lee
WorkingDirectory=/opt/kanban/看板正式程序
# 包装脚本内自带 42 重启 / 回滚 / 5 次失败停下；systemd Restart 作双保险
ExecStart=/bin/bash /opt/kanban/看板正式程序/deploy/linux/start_with_rollback.sh
Restart=always
RestartSec=3
# 日志走 journald：journalctl -u kanban -n 200 --no-pager
StandardOutput=journal
StandardError=journal
# 环境：UTF-8 兜底（中文路径/日志）
Environment=LANG=C.UTF-8
Environment=LC_ALL=C.UTF-8
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
# nginx 反代生产：后端只绑回环、不挂 static（由 nginx :80 伺服）
Environment=KANBAN_SERVER_HOST=127.0.0.1
Environment=KANBAN_SERVE_STATIC=0
Environment=KANBAN_SERVER_PORT=8018
# 离线抓数开关可由本地配置控制；此处不强行设 KANBAN_OFFLINE
# 任务书64·D8：进程沙箱（ReadWritePaths 须含程序/数据/日志/备份/归档 + 用户运行时）
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/kanban/看板正式程序 /opt/kanban/看板正式程序/数据 /opt/kanban/看板正式程序/数据/日志 /opt/kanban/看板正式程序/数据/备份 /opt/kanban/看板正式程序/数据/年度归档 /run/user/1000 /home/lee

[Install]
WantedBy=multi-user.target
